# src/systemd/sitedeploy.service
[Unit]
Description=Auto-run and deploy sites

[Service]
ExecStart=/usr/bin/sitedeploy
Restart=always

[Install]
WantedBy=multi-user.target
